# docker-mariadb/supervisord.conf
[supervisord]
nodaemon=true
 
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
  
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autorestart=true

[program:mariadb]
command=/usr/bin/mysqld_safe
user=mysql
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/mysql/mysqld.err.log
stdout_logfile=/var/log/mysql/mysqld.out.log
environment=MYSQLD_OPTS='--defaults-file=/etc/mysql/my.cnf'


